# supervisord.ini
[supervisord]
nodaemon=true

[program:main]
command=/home/web/main          ; this needs to be set as path to where compiled api app is
autostart=true                  ; start at supervisord start
autorestart=true                ; auto-restart on crash
stderr_logfile=/var/log/main.err.log
stdout_logfile=/var/log/main.out.log